# www/cron/domobot.cron
#!/bin/bash
gad &
php /var/www/html/domobot/cron/dapi.php &